# data/xdg-desktop-portal-faraway.service.in
[Unit]
Description=AppChooser portal that opens the default app

[Service]
Type=dbus
BusName=org.freedesktop.impl.portal.desktop.faraway
ExecStart=@LIBEXECDIR@/xdg-desktop-portal-faraway
